# systemd/ins-notify.service
[Unit]
Description=instantCLI notification history monitor

[Service]
Type=simple
ExecStart=/usr/bin/ins notify daemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
